# greet-sign.service
[Unit]
Description=LED greeter sign

[Service]
ExecStart=/usr/local/bin/greet-sign

[Install]
WantedBy=multi-user.target
