# greet-sign.service
[Unit]
Description=LED greeter sign

[Service]
ExecStart=/home/pi/ledsign.sh

[Install]
WantedBy=multi-user.target
